# mcp-manager.service
[Unit]
Description=MCP Manager - MCP Server Management Dashboard
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mcp-manager --port 9847 --config /etc/mcp-manager/config.json
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
